# config/systemd/ideal-magic-backup.timer
[Unit]
Description=Daily Ideal Magic PostgreSQL backups

[Timer]
# 03:30 UTC daily. The Scryfall corpus refresh runs at 10:30 UTC, so a 03:30 UTC
# backup window keeps the dump out of the heaviest write window.
OnCalendar=*-*-* 03:30:00 UTC
Persistent=true
RandomizedDelaySec=5min
Unit=ideal-magic-backup.service

[Install]
WantedBy=timers.target
